# docs/systemd/teelo-update.timer
[Unit]
Description=Run Teelo current events update on schedule

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
